# .crontab
@reboot cd /home/orangepi/laundry-alert && make api/prod